# services/fax/docker/core/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:asterisk]
command=/usr/sbin/asterisk -f -U asterisk -G asterisk -vvvg -c
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/asterisk/asterisk_stdout.log
stderr_logfile=/var/log/asterisk/asterisk_stderr.log
startretries=3
startsecs=5

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:iaxmodem]
command=/usr/bin/iaxmodem ttyIAX0
autostart=true
autorestart=true
priority=40
startsecs=3
stdout_logfile=/var/log/asterisk/iaxmodem_stdout.log
stderr_logfile=/var/log/asterisk/iaxmodem_stderr.log

[program:hfaxd]
command=/usr/sbin/hfaxd -d -i 4559
autostart=true
autorestart=true
priority=50
startsecs=3
stdout_logfile=/var/log/asterisk/hfaxd_stdout.log
stderr_logfile=/var/log/asterisk/hfaxd_stderr.log

[program:faxq]
command=/usr/sbin/faxq -D
autostart=true
autorestart=true
priority=50
startsecs=3
stdout_logfile=/var/log/asterisk/faxq_stdout.log
stderr_logfile=/var/log/asterisk/faxq_stderr.log

[program:faxgetty]
command=/usr/sbin/faxgetty ttyIAX0
autostart=true
autorestart=true
priority=60
startsecs=5
startretries=10
depends_on=iaxmodem
stdout_logfile=/var/log/asterisk/faxgetty_stdout.log
stderr_logfile=/var/log/asterisk/faxgetty_stderr.log

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pm2-fastagi]
command=pm2-runtime start /var/www/html/admin/modules/core/node/fastagi-server.js --name fastagi
autostart=true
autorestart=true
priority=70
user=asterisk
environment=HOME="/var/lib/asterisk",PM2_HOME="/var/lib/asterisk/.pm2"
stdout_logfile=/var/log/asterisk/pm2_stdout.log
stderr_logfile=/var/log/asterisk/pm2_stderr.log
